# modules/feh/XDG_CONFIG_HOME.feh/systemd/user/feh-wallpaper.path
[Path]
Unit=daemon-reload.service
PathChanged=%h/.config/environment.d/wallpaper.conf

[Install]
WantedBy=multi-user.target
