# modules/feh/XDG_CONFIG_HOME.feh/systemd/user/feh-fehbg.path
[Path]
Unit=daemon-reload.service
PathChanged=%h/.config/feh/.fehbg

[Install]
WantedBy=multi-user.target
