# modules/entr/XDG_CONFIG_HOME.entr/systemd/user/entr-watch-watcher.path
[Path]
Unit=entr-watch-watcher.service
PathChanged=%h/.config/entr/watchers

[Install]
WantedBy=default.target
