# modules/entr/XDG_CONFIG_HOME.entr/systemd/user/entr-watch-watcher.path
# deprecated
[Path]
Unit=entr-watch-watcher.service
PathChanged=%E/entr/watchers

[Install]
WantedBy=watcher.target
