# modules/xorg/.xorg/.config/systemd/user/xsession.target
[Unit]
Description=Xsession running
BindsTo=graphical-session.target
