# modules/xdg-services/XDG_CONFIG_HOME.xdg-services/systemd/user/xdg-runtime-reset.service
[Unit]
Description=Entr watch autorestarter

[Service]
Type=oneshot
ExecStart=%h/.local/bin/cleanruntime

[Install]
WantedBy=default.target
