# modules/sxhkd/XDG_CONFIG_HOME.sxhkd/systemd/user/sxhkd.path
# deprecated
[Path]
PathModified=%E/sxhkd/sxhkdrc
PathModified=%E/sxhkd/sxhkdrc.d

[Install]
WantedBy=default.target
