# modules/polybar/XDG_CONFIG_HOME.polybar/systemd/user/polybar.path
[Path]
Unit=polybar.service
PathChanged=%h/.config/polybar/polybar.sh

[Install]
WantedBy=default.target
