# modules/polybar/XDG_CONFIG_HOME.polybar/systemd/user/polybar-config.path
[Path]
Unit=polybar.service
PathChanged=%h/.config/polybar/config

[Install]
WantedBy=default.target
