# modules/tty-screensaver/resources/getty@tty1.service.d/override-bottom.conf
[Service]
Type=idle
ExecStart=
ExecStart=-btm --color gruvbox-light --dot_marker --group
StandardInput=tty
StandardOutput=tty
UnsetEnvironment=
