# modules/tty-screensaver/resources/getty@tty1.service.d/override-pipes.conf
[Service]
Type=idle
ExecStart=
ExecStart=-ss_pipes
StandardInput=tty
StandardOutput=tty
UnsetEnvironment=
