# modules/tty-screensaver/resources/getty@tty1.service.d/override-htop.conf
[Service]
Type=idle
ExecStart=
ExecStart=-/usr/bin/htop
StandardInput=tty
StandardOutput=tty
UnsetEnvironment=
