# modules/systemd/XDG_CONFIG_HOME.systemd/systemd/user/environment-watcher.target
[Unit]
Description=Listen environment changes
Wants=notify-restart@%n.service

[Install]
WantedBy=environment-watcher.service
