# modules/systemd/XDG_CONFIG_HOME.systemd/systemd/user/notify-restart@.service
[Unit]
Description=Notify about restarts

[Service]
Type=oneshot
ExecStart=notify-send --urgency=normal '%i restarted'

[Install]
WantedBy=default.target
