# modules/systemd/XDG_CONFIG_HOME.systemd/systemd/user/notify-low@.service
[Unit]
Description=Notify about low priority events

[Service]
Type=oneshot
ExecStart=notify-send --urgency=low '%i'

[Install]
WantedBy=default.target
