# modules/systemd/XDG_CONFIG_HOME.systemd/systemd/user/notify-normal@.service
[Unit]
Description=Notify about events

[Service]
Type=oneshot
ExecStart=notify-send --urgency=normal '%i'

[Install]
WantedBy=default.target
