# modules/systemd/XDG_CONFIG_HOME.systemd/systemd/user/user-services.path
# deprecated, see: user-service-watcher.service
[Path]
Unit=daemon-reload.service
PathModified=%E/systemd/user

[Install]
WantedBy=default.target
