# modules/systemd/XDG_CONFIG_HOME.systemd/systemd/user/user-services.path
[Path]
Unit=daemon-reload.service
PathChanged=%h/.config/systemd/user

[Install]
WantedBy=default.target
