# modules/sxhkd/.sxhkd/.config/systemd/user/sxhkd-watcher.path
[Path]
Unit=sxhkd-watcher.service
PathChanged=%h/.config/sxhkd/sxhkdrc.d
# PathModified=%h/.config/sxhkd/sxhkdrc

[Install]
WantedBy=default.target
