# modules/fehbg/XDG_CONFIG_HOME.$systemd.fehbg/systemd/user/feh-wallpaper.path
[Path]
Unit=daemon-reload.service
PathChanged=%E/environment.d/wallpaper.conf

[Install]
WantedBy=multi-user.target
